# backend/backend.service
[Unit]
Description=FastAPI backend for alberto-resume
After=network.target

[Service]
User=webdev
Group=webdev
WorkingDirectory=/home/webdev/projects/alberto-resume/backend
ExecStart=/home/webdev/projects/alberto-resume/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target